# services/unifi/run
#!/bin/bash

echo "Starting Unifi Controller."
/etc/init.d/unifi start
